# systemd/watchtower-oneshot.service
[Unit]
Description=Watchtower one-time update run
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker run --rm \
    -v /var/run/docker.sock:/var/run/docker.sock \
    containrrr/watchtower:latest \
    --run-once \
    --cleanup \
    --debug
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
